# deploy/cwatlas-collector.service
[Unit]
Description=CWAtlas CW collector (Web-888 -> MorseBase corpus)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dnelms
WorkingDirectory=/home/dnelms/cwatlas/collector
Environment=PYTHONUNBUFFERED=1
# Station details — SDR address, antenna location, data dir — live in
# config.toml, which is NOT tracked. Copy config.example.toml and edit it.
Environment=CWATLAS_CONFIG=/home/dnelms/cwatlas/collector/config.toml
# --no-mcp: stdio MCP needs a client on stdin; under systemd that's /dev/null.
# Attach an agent by running the collector manually, or wait for the
# streamable-HTTP transport milestone.
ExecStart=/home/dnelms/cwatlas/collector/.venv/bin/python -m cwatlas_mcp.runtime --no-mcp
# SIGTERM triggers the collector's graceful shutdown (workers finalize
# in-flight SigMF files + catalog rows); give it time before SIGKILL
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
