# deploy/cwatlas-dash.service
[Unit]
Description=CWAtlas status dashboard (read-only Flask sidecar, port 8828)
# Independent of cwatlas-collector by design (Model B): the dashboard must
# render (and show RED) precisely when the collector is down.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dnelms
WorkingDirectory=/home/dnelms/cwatlas/collector
Environment=PYTHONUNBUFFERED=1
# Same config.toml the collector reads, so both agree about where the station
# is. Not tracked — copy config.example.toml and edit it.
Environment=CWATLAS_CONFIG=/home/dnelms/cwatlas/collector/config.toml
ExecStart=/home/dnelms/cwatlas/collector/.venv/bin/python -m cwatlas_dash --host 0.0.0.0 --port 8828
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
