# deploy/cwatlas-dash.service
[Unit]
Description=CWAtlas status dashboard (read-only Flask sidecar, port 8828)
# Independent of cwatlas-collector by design (Model B): the dashboard must
# render (and show RED) precisely when the collector is down.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dnelms
WorkingDirectory=/home/dnelms/cwatlas/collector
Environment=PYTHONUNBUFFERED=1
Environment=CWATLAS_SDR_HOST=192.168.2.46
Environment=CWATLAS_DATA_DIR=/mnt/md0/cwatlas/data
Environment=CWATLAS_LAT=33.427
Environment=CWATLAS_LON=-82.208
ExecStart=/home/dnelms/cwatlas/collector/.venv/bin/python -m cwatlas_dash --host 0.0.0.0 --port 8828
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
